# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log

[program:nofx-backend]
command=/app/nofx
directory=/app
autostart=true
autorestart=true
startsecs=1
environment=TZ="%(ENV_NOFX_TIMEZONE)s",NOFX_ADMIN_PASSWORD="%(ENV_NOFX_ADMIN_PASSWORD)s"
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr




